# supervisord.conf
; supervisord.conf — AI Employee (Production Linux Server)
;
; supervisord is a battle-tested UNIX process manager. Unlike PM2 it has no Node.js
; dependency and is the standard choice for Python services on production servers.
;
; Install:
;   pip install supervisor        # or: apt install supervisor
;
; Start supervisord:
;   supervisord -c supervisord.conf
;
; Control:
;   supervisorctl -c supervisord.conf status
;   supervisorctl -c supervisord.conf restart orchestrator
;   supervisorctl -c supervisord.conf tail -f orchestrator
;   supervisorctl -c supervisord.conf stop all
;
; To run as a system service, copy to /etc/supervisor/conf.d/ai-employee.conf
; and run: sudo supervisorctl reread && sudo supervisorctl update

[unix_http_server]
file=/tmp/supervisor-ai-employee.sock
chmod=0700

[supervisord]
logfile=/mnt/d/Hackathon-00/Ai-Employee/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/tmp/supervisord-ai-employee.pid
nodaemon=false
silent=false
user=%(ENV_USER)s

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor-ai-employee.sock

; ── Master Orchestrator ───────────────────────────────────────────────────────

[program:orchestrator]
command=%(ENV_HOME)s/.local/bin/uv run python orchestrator.py --no-gmail --no-linkedin --no-social
directory=/mnt/d/Hackathon-00/Ai-Employee
autostart=true
autorestart=true
startsecs=10                    ; must stay up 10s to be considered started
startretries=3
stopwaitsecs=10
stdout_logfile=/mnt/d/Hackathon-00/Ai-Employee/logs/orchestrator.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stderr_logfile=/mnt/d/Hackathon-00/Ai-Employee/logs/orchestrator-error.log
stderr_logfile_maxbytes=2MB
environment=VAULT_PATH="/mnt/d/Hackathon-00/Ai-Employee/AI_Employee_Vault",DRY_RUN="false"
priority=100                    ; start first

; ── Scheduler ────────────────────────────────────────────────────────────────

[program:scheduler]
command=%(ENV_HOME)s/.local/bin/uv run scheduler
directory=/mnt/d/Hackathon-00/Ai-Employee
autostart=true
autorestart=true
startsecs=30
startretries=3
stdout_logfile=/mnt/d/Hackathon-00/Ai-Employee/logs/scheduler.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stderr_logfile=/mnt/d/Hackathon-00/Ai-Employee/logs/scheduler-error.log
stderr_logfile_maxbytes=2MB
environment=VAULT_PATH="/mnt/d/Hackathon-00/Ai-Employee/AI_Employee_Vault",DRY_RUN="false"
priority=200                    ; start after orchestrator

; ── Watchdog ─────────────────────────────────────────────────────────────────

[program:watchdog]
command=%(ENV_HOME)s/.local/bin/uv run watchdog-service
directory=/mnt/d/Hackathon-00/Ai-Employee
autostart=true
autorestart=true
startsecs=5
startretries=10
stdout_logfile=/mnt/d/Hackathon-00/Ai-Employee/logs/watchdog.log
stdout_logfile_maxbytes=2MB
stdout_logfile_backups=2
stderr_logfile=/mnt/d/Hackathon-00/Ai-Employee/logs/watchdog-error.log
stderr_logfile_maxbytes=1MB
environment=VAULT_PATH="/mnt/d/Hackathon-00/Ai-Employee/AI_Employee_Vault"
priority=300                    ; start last (after orchestrator is up)

; ── Gmail Watcher (uncomment after credentials configured) ───────────────────

;[program:gmail-watcher]
;command=%(ENV_HOME)s/.local/bin/uv run gmail-watcher
;directory=/mnt/d/Hackathon-00/Ai-Employee
;autostart=false
;autorestart=true
;startsecs=60
;startretries=3
;stdout_logfile=/mnt/d/Hackathon-00/Ai-Employee/logs/gmail-watcher.log
;stdout_logfile_maxbytes=5MB
;stdout_logfile_backups=3
;stderr_logfile=/mnt/d/Hackathon-00/Ai-Employee/logs/gmail-watcher-error.log
;stderr_logfile_maxbytes=2MB
;environment=VAULT_PATH="/mnt/d/Hackathon-00/Ai-Employee/AI_Employee_Vault"
;priority=400

; ── Process groups ───────────────────────────────────────────────────────────

[group:ai-employee-core]
programs=orchestrator,scheduler,watchdog

;[group:ai-employee-watchers]
;programs=gmail-watcher
